# teletraan/conf/snappass/testjob.conf
description "A test job file for experimenting with Upstart"
author "Your Name"
start on runlevel [2345]
exec echo Nick's Job ran at  `date` >> /var/log/testjob.log